# debian/postinst
#!/bin/bash

systemctl enable set-irq-affinity
