# etc/udev/rules.d/99-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{HDMI}=="1", RUN+="/usr/bin/hdmi-change"
