# etc/udev/rules.d/99-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/hdmi-change"
